# cron/crontab
0 0 * * * $HOME/dotfiles/cron/rm_dwl
0 15 * * * $HOME/dotfiles/cmd/maintenance
